# ansible/roles/mailhog/templates/mailhog_supervisor.conf.j2
[program:mailhog]
command={{ mailhog_install_dir }}/mailhog
autostart=true
autorestart=true
stdout_logfile=/var/log/mailhog.log
redirect_stderr=true